# etc/services.d/logpipe/run
#!/bin/sh
mkfifo -m 666 /run/logpipe
cat <> /run/logpipe 1>&2
