# etc/services.d/aphlict/run
#!/usr/bin/execlineb -P
s6-setuidgid nobody /phabricator/bin/aphlict debug --config /etc/aphlict.json
