# etc/services.d/php-fpm/run
#!/usr/bin/execlineb -P
foreground {
    s6-svc -uwu /run/s6/services/logpipe
}
foreground {
    s6-mkdir -p /run/php
}
/usr/sbin/php-fpm7 -FO
